# backend/jazzypop-quiz-generator.service
[Unit]
Description=JazzyPop Quiz Set Generator
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/jazzypop-backend
Environment="PATH=/home/ubuntu/jazzypop-backend/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ubuntu/jazzypop-backend/venv/bin/python /home/ubuntu/jazzypop-backend/quiz_set_generator.py
Restart=always
RestartSec=30
StandardOutput=append:/home/ubuntu/jazzypop-backend/logs/quiz_set_generator.log
StandardError=append:/home/ubuntu/jazzypop-backend/logs/quiz_set_generator_error.log

# Restart conditions
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target