# jazzypop-backup-pre-economy-20250703_160007/backend/jazzypop-backend.service
[Unit]
Description=JazzyPop Backend API
After=network.target postgresql.service redis.service
Wants=postgresql.service redis.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/jazzypop-backend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/ubuntu/jazzypop-backend/main.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/jazzypop-backend.log
StandardError=append:/var/log/jazzypop-backend.error.log

# Restart conditions
StartLimitIntervalSec=300
StartLimitBurst=5

[Install]
WantedBy=multi-user.target