# jazzypop-backup-pre-economy-20250703_160007/backend/jazzypop-monitor.service
[Unit]
Description=JazzyPop System Monitor
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/jazzypop-backend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
ExecStart=/usr/bin/python3 /home/ubuntu/jazzypop-backend/system_monitor.py
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target